# kombinace/lib/systemd/system/serverstart.service
[Unit]
Description=Startup service
After=network-online.service

[Service]
ExecStart=node /var/www/server/server.js
User=root

[Install]
WantedBy=multi-user.target
