# televize/lib/systemd/system/startup.service
[Unit]
Description=Startup service

[Service]
ExecStart=/usr/bin/startup
User=pi

[Install]
WantedBy=multi-user.target
