# kombinace/etc/crontab
15 1   * * *  root    systemctl restart startup
