# ciselnik/lib/systemd/system/donoff.service
[Unit]
Description=Monitor on/off service

[Service]
ExecStart=/usr/bin/donoff
User=pi

[Install]
WantedBy=multi-user.target
